# crontab
# mcd-donalds — pipeline ETL diaria
# Executa todo dia as 08:00 (fuso do container: America/Sao_Paulo).
#
# Campos: minuto hora dia-do-mes mes dia-da-semana  (* * * = todo dia)
#
# O job carrega o ambiente do container (. /app/container_env.sh) porque o cron
# NAO herda as variaveis do processo do container — sem isso, PG_HOST/PG_USER e
# afins nao chegariam ao pipeline. O entrypoint gera esse arquivo no start.
0 8 * * * . /app/container_env.sh; cd /app && mcd-donalds >> /app/dados/logs/cron.log 2>&1
